# apache/files/etc/systemd/system/httpd.service.d/override.conf
[Service]
Restart=on-failure
TimeoutStopSec=3
