# spark/templates/etc/profile.d/spark.sh.j2
export PATH=$PATH:{{ spark_root_dir }}/spark/bin:{{ spark_root_dir }}/spark/sbin
export SPARK_HOME={{ spark_root_dir }}/spark