# hadoop/templates/etc/profile.d/hadoop.sh.j2
export PATH=$PATH:{{ hadoop_root_dir }}/hadoop/bin:{{ hadoop_root_dir }}/hadoop/sbin
export HADOOP_HOME={{ hadoop_root_dir }}/hadoop