# rabbitmq/templates/etc/systemd/system/rabbitmq-server.service.d/limits.conf
[Service]
LimitNOFILE={{ rabbitmq_limit_nofiles }}
