# mariadb-galera-cluster/templates/etc/systemd/system/mariadb.service.d/override.conf.j2
[Service]
TimeoutSec={{ mariabd_systemd_timeoutsec }}
