# genieacs/templates/etc/systemd/system/genieacs-nbi.service.j2
[Unit]
Description=Genieacs NBI server
After=network.target

[Service]
Restart=on-failure
User={{genieacs_user}}
Group={{genieacs_group}}
ExecStart=/usr/local/bin/node {{ genieacs_install_dir }}/bin/genieacs-nbi
EnvironmentFile=/etc/default/genieacs

[Install]
WantedBy=multi-user.target
