# prosody/templates/etc/systemd/system/prosody.service.d/override.conf
[Service]
Type=simple
User=prosody
Group=prosody
